# dist/DEBIAN/postinst
#!/bin/bash

pip install -r requirements.txt
systemctl enable weather-station
systemctl start weather-station