# dist/DEBIAN/prerm
#!/bin/bash

systemctl stop weather-station
systemctl disable weather-station
